# nvme0-temp/nvme0-temp.service
[Unit]
Description=NVMe temperature monitor daemon

[Service]
ExecStart=/usr/bin/nvme0-temp

[Install]
WantedBy=multi-user.target
